# service/treinenapp-watcher.path
[Path]
PathModified=/share/Git/Treintafel

[Install]
WantedBy=multi-user.target